# .config/systemd/user/hyprsunset.service
[Unit]
Description=Change Screen Temperature

[Service]
Type=simple
ExecStart=/home/allisnull/.dotfiles/scripts/hyprsunset.sh
