# .config/systemd/user/hyprsunset.timer
[Unit]
Description=HyprSunset Hourly Check

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
