# .config/systemd/user/low-battery.service
[Unit]
Description=Low Battery Warning

[Service]
ExecStart=/home/allisnull/.dotfiles/scripts/low_battery.sh
